# adguard-unbound-redis/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: AdGuard Home with Unbound & Redis
# Generates nginx configuration for ingress
# ==============================================================================

declare ingress_interface
declare ingress_port
declare ingress_entry

# Only generate nginx ingress config when ingress is enabled
ingress_port=$(bashio::addon.ingress_port 2>/dev/null || echo "")
if bashio::var.has_value "${ingress_port}"; then
    bashio::log.info "Generating nginx configuration for ingress..."

    ingress_interface=$(bashio::addon.ip_address)
    ingress_entry=$(bashio::addon.ingress_entry)

    bashio::var.json \
        interface "${ingress_interface}" \
        port "${ingress_port}" \
        ingress_entry "${ingress_entry}" \
        | tempio \
            -template /etc/nginx/templates/ingress.gtpl \
            -out /etc/nginx/servers/ingress.conf

    bashio::log.info "Nginx configuration generated."
else
    bashio::log.info "Ingress not enabled, skipping nginx configuration."
    # Create marker so svc-nginx knows to skip
    touch /tmp/.nginx_disabled
fi
